# newsboat/.config/systemd/user/newsboat-reload.timer
[Unit]
Description=Timer to reload all feeds in newsboat

[Timer]
OnBootSec=2m

[Install]
WantedBy=timers.target
