# newsboat/.config/systemd/user/newsboat-reload.service
[Unit]
Description=Reload all feeds in newsboat

[Service]
ExecStart=newsboat -x reload
